# docker-images/ubuntu-target/supervisord.conf
[supervisord]
nodaemon=true

[program:npm]
command=npm start

[program:tcpdump]
command=tcpdump -n -C 10 -w /log/target_packets.pcap
